# rootfs/etc/s6-overlay/s6-rc.d/fcgiwrap/up
if { spawn-fcgi -s /var/run/fcgiwrap/fcgiwrap.sock -- /usr/bin/fcgiwrap }
chown nginx:www-data /var/run/fcgiwrap/fcgiwrap.sock
