# rootfs/etc/s6-overlay/s6-rc.d/fcgiwrap/up
if { spawn-fcgi -s /var/run/fcgiwrap/fcgiwrap.sock -u nginx -U nginx -- /usr/bin/fcgiwrap }
