# rootfs/etc/s6-overlay/s6-rc.d/fcgiwrap/up
if { spawn-fcgi -s /var/run/fcgiwrap/fcgiwrap.sock -u fcgiwrap -g www-data -- /usr/bin/fcgiwrap }
chown nginx:nginx /var/run/fcgiwrap/fcgiwrap.sock
